# deploy/linkscribe-certbot-renew.service
[Unit]
Description=Renew LinkScribe TLS certificate
After=network-online.target nginx.service

[Service]
Type=oneshot
ExecStart=/opt/certbot/bin/certbot renew --quiet --deploy-hook "systemctl reload nginx"
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
