# deploy/linkscribe-api.service
[Unit]
Description=LinkScribe API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=linkscribe
Group=linkscribe
WorkingDirectory=/opt/linkscribe/app
EnvironmentFile=/etc/linkscribe/linkscribe.env
ExecStart=/opt/linkscribe/venv/bin/uvicorn linkscribe.api:app --host 127.0.0.1 --port 8080 --workers 1 --no-access-log
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/linkscribe /var/log/linkscribe
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
