# deploy/linkscribe-certbot-renew.timer
[Unit]
Description=Twice-daily LinkScribe TLS renewal check

[Timer]
OnCalendar=*-*-* 03,15:17:00
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
